# systemd/unbound-blacklist.timer
[Unit]
Description=Run unbound blacklist update every 4 hours

[Timer]
OnCalendar=0/4:00:00
Persistent=true

[Install]
WantedBy=timers.target
